# pizero-openclaw.service
[Unit]
Description=Pi Zero OpenClaw Voice Terminal
After=network-online.target
Wants=network-online.target

[Service]
# ── Adjust these to match your setup ──────────────────────────────────────
# User and group to run as (set to your actual username, e.g. rocky)
User=pi
Group=pi
# Path to the project (set to where you cloned pizero-openclaw)
WorkingDirectory=/home/pi/pizero-openclaw
# Always go through the wrapper so .env parsing stays consistent and remote
# loopback/SSH-tunnel setups behave the same way in shell and systemd.
ExecStart=/home/pi/pizero-openclaw/run-openclaw.sh
# ───────────────────────────────────────────────────────────────────────────

# 自动重启：崩溃后重启 / 系统唤醒后重启
Restart=always
RestartSec=5

# 防止重启循环：如果 10 秒内重启超过 3 次，则停止重启（硬件故障时避免刷日志）
StartLimitBurst=3
StartLimitIntervalSec=10

# 日志
StandardOutput=journal
StandardError=journal

# 开机自启
[Install]
WantedBy=multi-user.target
